# pypi-server/pypi_supervisord.ini
[supervisord]
nodaemon=true

[program:pypi_server]
command=/usr/bin/gunicorn -w 1 -b 0.0.0.0:8080 --pythonpath='/usr/lib/python3.5/dist-packages,/usr/python3/dist-packages,/usr/lib/python3.5' pypi-server
directory=/~/packages
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/opt/supervisor/pypi_server.log
